# services/blckswan-watchdog.service
[Unit]
Description=BLCKSWAN watchdog

[Service]
Type=oneshot
ExecStart=%h/blckswan/services/blckswan-watchdog.sh
